# tools/systemd/nonraid-parity-check.timer
[Unit]
Description=NonRAID quarterly parity check timer

[Timer]
OnCalendar=quarterly
Persistent=true
RandomizedDelaySec=3600
Unit=nonraid-parity-check.service

[Install]
WantedBy=timers.target
